# profile.d/editor.sh
export VISUAL=/usr/bin/vim
export EDITOR=$VISUAL

