# profile.d/history.sh
#!/usr/bin/env bash

export HISTTIMEFORMAT="[%h/%d/%Y - %H:%M:%S]"
export HISTSIZE=2000;
export HISTFILESIZE=2000;
