# etc/profile.d/browser.sh
#!/usr/bin/env bash
export BROWSER=google-chrome
