# etc/profile.d/vim_alias.sh
#!/usr/bin/env bash

#alias vim='vimx'
alias vim='gvim -v'
